# provisioning/etc/systemd/system/rb_sync.service
[Service]
ExecStart=/provisioning/bin/run_sync.sh
