# provisioning/etc/systemd/system/rb_watchdog.service
[Service]
ExecStart=/provisioning/bin/run_watchdog.sh
